# srv/primary/dcron/cron.d/backup
# MIN HOUR DAY MONTH DAYOFWEEK   COMMAND
* * * * * ID=backup FREQ=1d backup.sh --allow-source-mismatch
