# srv/primary/dcron/cron.d/health-check
# MIN HOUR DAY MONTH DAYOFWEEK   COMMAND
* * * * * touch /health-check
