# srv/secondary/dcron/cron.d/djbdns
# MIN HOUR DAY MONTH DAYOFWEEK   COMMAND
* * * * * docker exec tinydns-secondary make >/dev/null
